# iso/config/includes.chroot/usr/lib/systemd/user/aios.target
[Unit]
Description=AIOS Application Suite
Wants=aios-agent.service aios-chat.service aios-dock.service aios-confirm.service

[Install]
WantedBy=default.target
