# iso/config/includes.chroot/usr/lib/systemd/user/aios-chat.service
[Unit]
Description=AIOS Chat UI
After=aios-agent.service
Requires=aios-agent.service

[Service]
Type=simple
ExecStart=/usr/local/bin/aios-chat
Restart=on-failure
RestartSec=2
Environment=RUST_LOG=info

[Install]
WantedBy=default.target
